# scripts/power-monitor-shutdown.service
[Unit]
Description=Log power state for power-monitor
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo "{\"ts\":$(date +%%s),\"action\":\"pre\",\"what\":\"shutdown\",\"sleep_action\":\"\"}" >> /var/lib/power-monitor/state-log.jsonl'

[Install]
WantedBy=shutdown.target reboot.target halt.target
